# etc/supervisord.d/rocketchat.ini
[program:rocketchat]
command=/opt/qnib/rocketchat/bin/start.sh
redirect_stdery=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
autostart=true
startsecs=5
exitcodes=0
