# etc/supervisord.d/rocketchat.ini
[program:rocketchat]
command=node main.js
directory=/opt/rocketchat/bundle/
redirect_stdery=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
autostart=true
startsecs=5
exitcodes=0
